# scripts/cookbooks/app/templates/default/app.conf.erb
description "gunicorn app"
author      "Edwin Knuth"


start on filesystem
stop on runlevel S

script

    export HOME="/usr/local/apps/geosurvey/server"
    export VENV="/usr/local/venv/geosurvey"
    cd $HOME
    exec sudo -u nginx $VENV/bin/gunicorn -b 0.0.0.0:8888 wsgi
end script
